# linux/systemd/flexgit.timer
[Unit]
Description=Run flexgit twice daily

[Timer]
OnCalendar=*-*-* 00,12:00:00
Persistent=true

[Install]
WantedBy=timers.target
